# midnight-ntwrk/genai-dev-platform/container-scripts/profile.d/compact.sh
# Midnight Compact compiler environment
export COMPACT_HOME="/usr/local/share/compact"
export COMPACT_DIRECTORY="/usr/local/share/compact-compilers"
export PATH="${COMPACT_HOME}:${PATH}"
